# .docker/app/services.d/horizon/run
#!/usr/bin/with-contenv bash

# Switch to web user
su webuser

#Bring up horizon
php /var/www/html/artisan horizon